# ollama-proxy.service
[Unit]
Description=Ollama proxy collector
After=network.target docker.service
Requires=docker.service

[Service]
User=tonymm81
Group=tonymm81
WorkingDirectory=/opt/ollama_proxy
ExecStart=/opt/ollama_proxy/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8080 --workers 1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
